# Keepass_Backup/99-keepass-backup.rules

ATTRS{idVendor}=="090c", ATTRS{idProduct}=="1000", RUN+="/bin/echo connected >> /tmp/keepass_backup/pipe"

